# overlay-rpiduino/lib/systemd/system/rpiduino.service
[Unit]
Description=rpiduino presentation

[Service]
ExecStart=/usr/bin/qmlscene /usr/share/rpiduino/slides.qml

[Install]
WantedBy=multi-user.target
